# roles/tars/templates/upstart/upstart.conf
description "{{ user }} {{ service }}"

start on started networking
stop on runlevel [016]

limit nofile 1000000 1000000

respawn
respawn limit 5 30

setuid {{ user }}
setgid {{ user }}

chdir /home/{{ user }}/{{ service }}
env NODE_ENV={{node_env}}
exec npm start
